# conf/supervisord-sample.conf
[unix_http_server]
file=/tmp/meduele_supervisor.sock   ; (the path to the socket file)

[supervisord]
logfile=/tmp/meduele_supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB       ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10          ; (num of main logfile rotation backups;default 10)
loglevel=info               ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/meduele_supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false              ; (start in foreground if true;default false)
minfds=1024                 ; (min. avail startup file descriptors;default 1024)
minprocs=200                ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/meduele_supervisor.sock ; use a unix:// URL  for a unix socket


[program:gunicorn]
command = /path/to/virtualenvs/meduele-lib/bin/gunicorn -c /path/to/gunicorn/conf.py run:app
directory = /path/to/project/meduele/serve
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/path/to/log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=20
